# deploy/certbot_renew.cron
# Check certificate renewal twice per day.
# Run as root: crontab -e  then paste this line.
17 */12 * * * certbot renew --quiet --deploy-hook "systemctl reload nginx"
